# bots/supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
silent=false
minfds=1024
minprocs=200

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

# [program:abraham]
# command=python /marsbots/bot.py abraham
# autostart=true
# autorestart=true

# [program:superfsmon-abraham]
# command=superfsmon /marsbots/bots/abraham abraham

[program:eden]
command=python /marsbots/bot.py eden_bot
autostart=true
autorestart=true

[program:superfsmon-eden]
command=superfsmon /marsbots/bots/abraham eden_bot
